# com.turtledome.BGChanger.plist.erb
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.turtledome.BGChanger</string>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>0</integer>
	</dict>
	<key>ProgramArguments</key>
	<array>
		<string>ruby</string>
		<string><%= File.expand_path(File.dirname(__FILE__)) %>/cbg.rb</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string><%= ENV['HOME'] %>/.rbenv/bin:/usr/local/bin:<%= ENV['HOME'] %>/.rbenv/shims:<%= ENV['PATH'] %></string>
	</dict>
</dict>
</plist>
